# home/azuwis/backup/crontab
0 16 * * * my-installed-package > ~/backup/installed-package
