# home/azuwis/backup/crontab
0 16 * * * /usr/local/bin/my-installed-package > ~/backup/installed-package
