# home/azuwis/backup/crontab
PATH=$PATH:/usr/local/bin
0 16 * * * my-installed-package > ~/backup/installed-package
