# home/azuwis/misc/crontab
0 16 * * * /usr/local/bin/my-installed-package > ~/misc/installed-package
